# roles/aws-s3-backup/files/aws-s3-backup.target
[Unit]
Description=AWS S3 backup target
StopWhenUnneeded=yes
