# aws-s3-backup/files/aws-s3-backup.timer
[Unit]
Description=Daily AWS S3 backup

[Timer]
OnCalendar=daily
Persistent=true
Unit=aws-s3-backup.target

[Install]
WantedBy=timers.target
