# goaccess/templates/goaccess.timer.j2
[Unit]
Description=GoAccess refresh timer

[Timer]
OnCalendar={{goaccess_timer_freq}}
Unit=goaccess.service

[Install]
WantedBy=timers.target
